# scripts/git-trending.crontab
45 2 * * *   /home/bob/bin/git-trending.sh
